# scripts/systemd/s3desk-lighthouse.timer
[Unit]
Description=Schedule Lighthouse CI for S3Desk

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
